# backend/supervisor.ini
[program:api_carol]
command=/root/.pyenv/versions/api.carol.orango.io/bin/gunicorn -k uvicorn.workers.UvicornWorker main:app --bind=0.0.0.0:8083 --workers=1
directory=/var/www/api.carol.orango.io
environment=env=PRODUCTION
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api.carol.err.log
stdout_logfile=/var/log/supervisor/api.carol.out.log